# startup/systemd/swupdate-progress.service
[Unit]
Description=swupdate progress service
Before=swupdate.service

[Service]
ExecStart=/usr/bin/swupdate-progress -r -w

[Install]
WantedBy=swupdate.service
